# ops/mise-backup.timer
[Unit]
Description=Mise nightly SQLite snapshot (02:30)

[Timer]
OnCalendar=*-*-* 02:30:00
Persistent=true
RandomizedDelaySec=300

[Install]
WantedBy=timers.target
